# _system/com.lexery.wiki-maintenance.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.lexery.wiki-maintenance</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__VAULT_ROOT__/_system/scripts/run-maintenance-launchd.sh</string>
    </array>
    <key>WorkingDirectory</key>
    <string>__VAULT_ROOT__</string>
    <key>RunAtLoad</key>
    <true/>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>10</integer>
    </dict>
    <key>ThrottleInterval</key>
    <integer>600</integer>
    <key>StandardOutPath</key>
    <string>__VAULT_ROOT__/_system/logs/launchd-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>__VAULT_ROOT__/_system/logs/launchd-stderr.log</string>
</dict>
</plist>
